# usr_data/lib/systemd/system/sddm.service.d/override.conf
[Unit]
Conflicts=kader-bootsplash.service
After=kader-bootsplash.service