# usr_data/lib/systemd/user/kader42-mode-daemon.service
[Unit]
Description=Kader Hybrid Mode Daemon
After=graphical-session.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/kader-mode-daemon.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target