# usr_data/lib/systemd/user/kader42-watcher.service
[Unit]
Description=Kader⁴² Tablet-Mode Watcher
# Waiting until the graphical session is ready
After=plasma-core.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kader42-watcher
# If the script crashes (e.g., D-Bus restart), restart it after 5 seconds
Restart=on-failure
RestartSec=5
# Prevents the process from getting “lost” in the session
Slice=session.slice

[Install]
# This ensures that the service starts as soon as KDE loads
WantedBy=graphical-session.target