# kader42-packages/kader-tools-pkg/kader42-gravity-sync.timer
[Unit]
Description=Run Kader42 Gravity Sync every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target