# releng/airootfs/etc_tmp/systemd/system/io-sensor-proxy.service.d/order.conf
[Unit]
After=plymouth-killer.service
Requires=plymouth-killer.service