# releng/airootfs/etc/xdg/autostart/kader42-watcher.desktop
[Desktop Entry]
Type=Application
Name=Kader42 Tablet Mode Watcher
Comment=Monitors tablet mode and switches layouts
Exec=/usr/bin/kader42-watcher
Icon=input-tablet
Terminal=false
Categories=System;
# Ensures that it only starts in KDE/Plasma (in case you ever try GNOME)
OnlyShowIn=KDE;
# Prevents the script from running multiple times
X-KDE-UniqueApplet=true
# --- The magic ingredient ---
# Phase 2 is “Applications.” This means Plasma is already loaded.
X-KDE-autostart-phase=2
# We tell KDE not to wait for the script to finish (since it's a loop)
X-KDE-autostart-after=panel