# etc_conf/skel/.config/autostart/show-welcome-screen.desktop
[Desktop Entry]
Type=Application
Name=Kader⁴² Welcome
Exec=/usr/bin/python /usr/bin/kader42-welcome
Icon=/usr/share/pixmaps/mello-with-tablet.png
Terminal=false
Categories=System;Installer;
X-KDE-Authorized-Link=true