# releng/airootfs/usr_tmp/lib/systemd/system/plasma-convertible-switcher.service
[Unit]
Description=Plasma Convertible Mode Switcher
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/plasma-convertible-switcher
Restart=always
RestartSec=5
Environment=XDG_RUNTIME_DIR=%t/run/user/%U
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/run/user/%U/bus

[Install]
WantedBy=default.target