# releng/airootfs/usr_tmp/lib/systemd/system/plymouth-killer.service
[Unit]
Description=Plymouth Killer Interceptor
DefaultDependencies=no
After=systemd-udev-settle.service
# WICHTIG: Wir nehmen display-manager.service hier raus, 
# damit wir nicht VOR ihm fertig sein MÜSSEN, sondern PARALLEL agieren.
Before=graphical.target
Conflicts=plymouth-quit.service

[Service]
Type=simple
# Wir entfernen das 'terminate' aus der Kette, da es oft den Buffer leert.
# Wir setzen ein winziges 'sleep' VOR den Kill, damit das Logo Zeit zum Rendern hat.
# ExecStart=/usr/bin/sh -c "/usr/bin/sleep 0.5; /usr/bin/plymouth quit --retain-splash || /usr/bin/killall -9 plymouthd"
ExecStart=/usr/bin/sh -c "\
  /usr/bin/plymouth display-update; \
  /usr/bin/sleep 0.1; \
  /usr/bin/plymouth quit --retain-splash; \
  /usr/bin/sleep 0.2; \
  /usr/bin/killall -q plymouthd || true"
# Der Post-Sleep gibt dem i915 Treiber die Millisekunden, die er braucht.
ExecStartPost=/usr/bin/sleep 0.1
SuccessExitStatus=0 1 143

[Install]
WantedBy=sysinit.target